# zcash.conf
description "zcash"
author "Ioannis Kokkinis"
start on runlevel [2345]
exec /usr/local/bin/zcash.sh
