# zcash.conf
description "zcash"
start on startup
task
exec ~/zcash.sh
